# pihole-stats.timer
[Unit]
Description=Run pihole-stats every 30 min

[Timer]
OnBootSec=15sec
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target